# run_krystof.cron
*/15 * * * * python3 run.py
*/20 * * * * python3 run.py -c config_krystof.yaml

